# chili/etc/skel/.config/autostart/icons.desktop
[Desktop Entry]
Exec=bigicons-desktop
Name=bigicons-desktop
Type=Application

